# ss/supervisor.d/kcptun.ini
[program:kcptun]
command=/opt/program/kcptun/server_linux_amd64 -t 127.0.0.1:53 -l :3307 --key 123456 --crypt aes-128 --mode fast2
process_name=kcptun
numprocs=1
stopsignal=INT
directory=/opt/program/kcptun
user=root
autorestart = true
redirect_stderr=true
stdout_logfile=/data/log/supervisor/kcptun.out
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stopasgroup=true
killasgroup=true
