# ss/supervisor.d/privoxy.ini
[program:privoxy]
command=/usr/sbin/privoxy  --no-daemon /opt/program/privoxy/etc/config
process_name=privoxy
numprocs=1
stopsignal=INT
directory=/opt/program/privoxy
user=root
autorestart = true
redirect_stderr=true
stdout_logfile=/data/log/supervisor/privoxy.out
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stopasgroup=true
killasgroup=true
