# ss/supervisor.d/ss.ini
[program:ss]
command=/opt/program/ss/shadowsocks-server
process_name=ss
numprocs=1
stopsignal=INT
directory=/opt/program/ss/
user=root
redirect_stderr=true
stdout_logfile=/data/log/supervisor/ss.out
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stopasgroup=true
killasgroup=true
autorestart = true
