# ss/supervisor.d/ssclient.ini
[program:ssclient]
command=/opt/program/ssclient/shadowsocks-local-linux64-1.1.5 -c /opt/program/ssclient/ssclient.json
process_name=ssclient
numprocs=1
stopsignal=INT
directory=/opt/program/ssclient
user=root
autorestart = true
redirect_stderr=true
stdout_logfile=/data/log/supervisor/ssclient.out
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stopasgroup=true
killasgroup=true
